# src/systemd/ps3remote.service
[Unit]
Description=Userspace PS3 Remote Driver
Wants=bluetooth.target

[Service]
ExecStart=/usr/bin/ps3remote

[Install]
WantedBy=multi-user.target
